# ansible/roles/5-cmdb-file-update/files/cmdb_update.service
[Unit]
Description=CMDB Refresh and AppDynamics Config Update

[Service]
Type=simple
ExecStart=/home/appdynamics/get_cmdb.sh

[Install]
WantedBy=multi-user.target
