# ansible/roles/6-motd-update/files/motd_update.service
[Unit]
Description=Refresh motd with details from A2RM

[Service]
Type=simple
ExecStart=/root/update_motd.sh

[Install]
WantedBy=multi-user.target
